# scripts/headless/ipr-provision.service
[Unit]
Description=IPR Keyboard Headless Wi-Fi Hotspot
After=network.target NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/ipr-provision.sh
User=root

[Install]
WantedBy=multi-user.target
